# supervisord.conf
[supervisord]
nodaemon=true

[program:laravel-queue-worker]
command=php artisan queue:work database --sleep=3 --tries=3
autostart=true
autorestart=true
stderr_logfile=/var/log/worker.err.log
stdout_logfile=/var/log/worker.out.log
user=root
directory=/var/www/app
